# ops/systemd/hermes-memory-rollup.timer
[Unit]
Description=Weekly memory roll-up

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3h

[Install]
WantedBy=timers.target
